# device/linux-gadget/systemd/ed25519-token-daemon.service
[Unit]
Description=Ed25519 token signing daemon
After=ed25519-token-gadget.service
Requires=ed25519-token-gadget.service

[Service]
Type=simple
User=edtoken
Group=edtoken
SupplementaryGroups=dialout
ExecStart=/usr/local/libexec/ed25519-token-daemon
Restart=always
RestartSec=1

NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/etc/ed25519-token-auth
ReadWritePaths=/run
DeviceAllow=/dev/ttyGS0 rw

[Install]
WantedBy=multi-user.target
